# monitoring/pihole-ha/systemd/pihole-sync.service
[Unit]
Description=Pi-hole Configuration Sync Service
After=network.target pihole-FTL.service
Wants=pihole-FTL.service

[Service]
Type=oneshot
ExecStart=/usr/local/bin/sync-pihole-config.sh --auto
StandardOutput=journal
StandardError=journal
SyslogIdentifier=pihole-sync

[Install]
WantedBy=multi-user.target
